# testdata/ignitions/systemd/bird.service
[Unit]
Description=bird

[Service]
ExecStart=/usr/bin/bird

[Install]
WantedBy=multi-user.target
