# testdata/test/systemd/bird.service
[Unit]
Description=bird

[Service]
ExecStart=/usr/bin/bird --server={{ MyURL }}

[Install]
WantedBy=multi-user.target
