# testdata/base/systemd/chronyd.service
[Unit]
Description=Chrony

[Service]
ExecStart=/usr/bin/chronyd

[Install]
WantedBy=multi-user.target
